# package/luci-app-autoupdate/root/etc/cron.d/autoupdate
# Hourly scheduler; autoupdate.sh enforces the LuCI interval setting.
17 * * * * /usr/bin/autoupdate.sh scheduled >> /var/log/autoupdate.log 2>&1
